# apps/macos/TalkieServer/launchd/com.talkie.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.talkie.bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>cd "${TALKIE_BRIDGE_DIR:-$HOME/dev/talkie/apps/macos/TalkieServer}" &amp;&amp; exec /opt/homebrew/bin/bun run src/server.ts</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/talkie-bridge.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/talkie-bridge.error.log</string>
</dict>
</plist>
