# macOS/TalkieBridge/launchd/com.talkie.bridge-dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.talkie.bridge-dev</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/bun</string>
        <string>run</string>
        <string>src/server.ts</string>
        <string>--local</string>
        <string>--port</string>
        <string>8767</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/arach/dev/talkie-action-state-machine/macOS/TalkieBridge</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/arach/Library/Logs/TalkieBridge/bridge-dev.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/arach/Library/Logs/TalkieBridge/bridge-dev.error.log</string>
</dict>
</plist>
