# packaging/macos/resources/user-agents/to.talkie.app.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>to.talkie.app.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Talkie.app/Contents/Library/LoginItems/TalkieSync.app/Contents/MacOS/TalkieSync</string>
    </array>
    <key>MachServices</key>
    <dict>
        <key>to.talkie.app.sync.xpc</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
