# packaging/macos/resources/user-agents/to.talkie.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>to.talkie.agent</string>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>to.talkie.app.mac</string>
        <string>to.talkie.agent</string>
    </array>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Talkie.app/Contents/Library/LoginItems/TalkieAgent.app/Contents/MacOS/TalkieAgent</string>
    </array>
    <key>MachServices</key>
    <dict>
        <key>to.talkie.agent.xpc</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
